# roles/prometheus/templates/prometheus.service.j2
[Unit]
Description=Prometheus Service
After=network.target

[Service]
User=root
ExecStart=/usr/local/bin/prometheus

[Install]
WantedBy=multi-user.target
